# .my_scripts/init/system/etc/systemd/system/clean-cache.service
[Unit]
Description=Clean cache

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Scc --noconfirm
