# .my_scripts/init/system/etc/systemd/system/clean-cache.service
[Unit]
Description=Clear Pacman Cache

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Scc
