# .my_scripts/init/system/etc/systemd/system/clean-cache.timer
[Unit]
Description=Run Pacman Cache Monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
