# .my_scripts/init/system/etc/systemd/system/denyhosts.timer
[Unit]
Description=Run denyhosts service after boot

[Timer]
OnBootSec=5min
Unit=denyhosts.service

[Install]
WantedBy=timers.target
