# .my_scripts/init/system/etc/systemd/system/mirrors.timer
[Unit]
Description=Run mirrors.service weekly

[Timer]
OnCalendar=weekly
OnBootSec=1min
Persistent=true

[Install]
WantedBy=timers.target
