# .my_scripts/init/system/etc/systemd/system/mirrors.timer
[Unit]
Description=Run mirrors.service weekly

[Timer]
OnBootSec=15min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target
