# .my_scripts/init/system/etc/systemd/system/clamav.daemon.service.d/override.conf
# Lowest priority

[Service]
Nice=19
CPUQuota=5%
CPUWeight=1
CPUSchedulingPolicy=idle
IOWeight=1
IOSchedulingClass=idle
IOSchedulingPriority=7
