# .my_scripts/init/system/etc/systemd/system/cups.service.d/delay.conf
# Delay by 1 minute to improve boot time

[Unit]
After=sysinit.target

[Timer]
OnBootSec=1min
