# .my_scripts/init/amd_oc/amd-overclock.service
[Unit]
Description=Amd overclock

[Service]
Type=simple
ExecStart=/usr/local/bin/oc.sh

[Install]
WantedBy=multi-user.target
